# config/init/heaverd.service
[Unit]
Description=Heaver daemon
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/heaverd

[Install]
WantedBy=multi-user.target
